# launchd/com.bmalan.kokoro-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bmalan.kokoro-tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>tool</string>
        <string>run</string>
        <string>--from</string>
        <string>kokoro-onnx</string>
        <string>--with</string>
        <string>sounddevice</string>
        <string>python</string>
        <string>/Users/bmalan/Projects/claude-tts/server/kokoro_server.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/kokoro-server.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kokoro-server.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bmalan</string>
    </dict>
</dict>
</plist>
